# ops/aws-testnet/runtime/systemd/scuro-prover-worker.service
[Unit]
Description=Scuro private prover worker
After=network-online.target scuro-anvil.service
Wants=network-online.target scuro-anvil.service

[Service]
Type=simple
EnvironmentFile=-/etc/scuro-testnet/scuro.env
EnvironmentFile=-/etc/scuro-testnet/bootstrap.env
EnvironmentFile=-/etc/scuro-testnet/runtime.env
WorkingDirectory=/opt/scuro-testnet/current
ExecStart=/opt/scuro-testnet/tools/bin/bun run /opt/scuro-testnet/current/ops/aws-testnet/src/worker.ts
Restart=always
RestartSec=3
StandardOutput=append:/var/log/scuro-testnet/prover-worker.log
StandardError=append:/var/log/scuro-testnet/prover-worker.log

[Install]
WantedBy=multi-user.target
